# mfpf-server/server/config/logs/rsyslog.conf
[program:startsyslog]
command=/bin/bash -c "service rsyslog start"