# mfpf-analytics/server/bin/liberty.conf
[program:liberty]
command=/opt/ibm/wlp/bin/liberty-run